# systemd/amoskys-eventbus.service
[Unit]
Description=AMOSKYS EventBus gRPC Server
After=network.target

[Service]
Type=simple
User=amoskys
Group=amoskys
WorkingDirectory=/opt/amoskys
Environment="PATH=/opt/amoskys/.venv/bin"
EnvironmentFile=/etc/amoskys/production.env
ExecStart=/opt/amoskys/.venv/bin/python src/amoskys/eventbus/server.py --overload off
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/amoskys/data /var/log/amoskys
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
