# deployments/kernel_audit/amoskys-kernel-audit-simple.service
[Unit]
Description=AMOSKYS KernelAudit Guard v2 - Syscall-Plane Threat Detection
After=auditd.service network.target
Requires=auditd.service
ConditionPathExists=/var/log/audit/audit.log

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Agent command - using venv Python
ExecStart=/home/ubuntu/amoskys-venv/bin/python3 /home/ubuntu/kernel_audit/run_agent_v2.py \
    --device-id=%H \
    --queue-path=/var/lib/amoskys/queues/kernel_audit \
    --audit-log=/var/log/audit/audit.log \
    --collection-interval=5 \
    --metrics-interval=60

# Restart policy
Restart=always
RestartSec=10s
StartLimitInterval=300s
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu/amoskys-src/src"
WorkingDirectory=/home/ubuntu/amoskys-src

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=amoskys-kernel-audit

[Install]
WantedBy=multi-user.target
