# deploy/agents/kernel_audit/audit_rules.conf
## AMOSKYS KernelAudit Agent - Audit Rules Configuration
##
## This file contains Linux audit rules optimized for KernelAuditGuardV2 probes.
## Install to: /etc/audit/rules.d/amoskys-kernel.rules
##
## After installation:
##   sudo augenrules --load
##   sudo systemctl restart auditd
##
## Performance Note: These rules generate high volume on busy systems.
## Monitor /var/log/audit/audit.log size and adjust as needed.

# Remove any existing rules (optional - uncomment if needed)
# -D

# Buffer and rate limit configuration
# Increase backlog to handle bursts
-b 8192

# Continue logging even if audit queue is full (vs. panic/halt)
-f 1

# =============================================================================
# Execution Monitoring (ExecveHighRiskProbe)
# =============================================================================

# Track all execve/execveat syscalls
# Used by: ExecveHighRiskProbe
-a always,exit -F arch=b64 -S execve -k amoskys_exec
-a always,exit -F arch=b64 -S execveat -k amoskys_exec

# Alternative: Only track execve from high-risk directories (lower volume)
# -a always,exit -F arch=b64 -S execve -F dir=/tmp -k amoskys_exec_tmp
# -a always,exit -F arch=b64 -S execve -F dir=/dev/shm -k amoskys_exec_shm
# -a always,exit -F arch=b64 -S execve -F dir=/var/tmp -k amoskys_exec_var

# =============================================================================
# Privilege Escalation Monitoring (PrivEscSyscallProbe)
# =============================================================================

# Track all privilege-changing syscalls
# Used by: PrivEscSyscallProbe
-a always,exit -F arch=b64 -S setuid -S seteuid -S setreuid -S setresuid -k amoskys_privesc
-a always,exit -F arch=b64 -S setgid -S setegid -S setregid -S setresgid -k amoskys_privesc
-a always,exit -F arch=b64 -S setfsuid -S setfsgid -k amoskys_privesc
-a always,exit -F arch=b64 -S capset -k amoskys_privesc

# =============================================================================
# Kernel Module Monitoring (KernelModuleLoadProbe)
# =============================================================================

# Track kernel module operations
# Used by: KernelModuleLoadProbe
-a always,exit -F arch=b64 -S init_module -S finit_module -k amoskys_module_load
-a always,exit -F arch=b64 -S delete_module -k amoskys_module_unload

# Watch for module loading utilities
-w /usr/sbin/insmod -p x -k amoskys_module_tools
-w /usr/sbin/rmmod -p x -k amoskys_module_tools
-w /usr/sbin/modprobe -p x -k amoskys_module_tools

# =============================================================================
# Process Injection Monitoring (PtraceAbuseProbe)
# =============================================================================

# Track ptrace and memory manipulation syscalls
# Used by: PtraceAbuseProbe
-a always,exit -F arch=b64 -S ptrace -k amoskys_ptrace
-a always,exit -F arch=b64 -S process_vm_readv -S process_vm_writev -k amoskys_proc_mem

# =============================================================================
# File Permission Tampering (FilePermissionTamperProbe)
# =============================================================================

# Track permission changes on /etc (high-value config directory)
# Used by: FilePermissionTamperProbe
-a always,exit -F arch=b64 -S chmod -S fchmod -S fchmodat -F dir=/etc -k amoskys_perm_etc
-a always,exit -F arch=b64 -S chown -S fchown -S lchown -S fchownat -F dir=/etc -k amoskys_perm_etc

# Watch specific sensitive files
-w /etc/passwd -p wa -k amoskys_passwd
-w /etc/shadow -p wa -k amoskys_shadow
-w /etc/group -p wa -k amoskys_group
-w /etc/gshadow -p wa -k amoskys_gshadow
-w /etc/sudoers -p wa -k amoskys_sudoers
-w /etc/sudoers.d/ -p wa -k amoskys_sudoers

# SSH configuration
-w /etc/ssh/sshd_config -p wa -k amoskys_sshd_config
-w /etc/ssh/ssh_host_rsa_key -p wa -k amoskys_ssh_host_key
-w /etc/ssh/ssh_host_ed25519_key -p wa -k amoskys_ssh_host_key

# =============================================================================
# Audit Tampering Detection (AuditTamperProbe)
# =============================================================================

# Watch audit configuration and logs
# Used by: AuditTamperProbe
-w /etc/audit/ -p wa -k amoskys_audit_config
-w /etc/audit/audit.rules -p wa -k amoskys_audit_rules
-w /etc/audit/auditd.conf -p wa -k amoskys_audit_conf
-w /var/log/audit/ -p wa -k amoskys_audit_logs
-w /var/log/audit/audit.log -p wa -k amoskys_audit_log

# Watch audit binaries
-w /usr/sbin/auditctl -p x -k amoskys_audit_tools
-w /usr/sbin/auditd -p x -k amoskys_audit_tools

# Track audit service management
-w /usr/lib/systemd/system/auditd.service -p wa -k amoskys_audit_service
-w /etc/systemd/system/auditd.service -p wa -k amoskys_audit_service

# =============================================================================
# Optional: Additional Security-Critical Paths
# =============================================================================

# PAM configuration (authentication)
-w /etc/pam.d/ -p wa -k amoskys_pam

# Cron jobs (persistence)
-w /etc/cron.d/ -p wa -k amoskys_cron
-w /etc/cron.daily/ -p wa -k amoskys_cron
-w /etc/cron.hourly/ -p wa -k amoskys_cron
-w /etc/cron.weekly/ -p wa -k amoskys_cron
-w /etc/cron.monthly/ -p wa -k amoskys_cron
-w /var/spool/cron/ -p wa -k amoskys_cron

# Systemd services (persistence)
-w /etc/systemd/system/ -p wa -k amoskys_systemd
-w /usr/lib/systemd/system/ -p wa -k amoskys_systemd

# =============================================================================
# Make Configuration Immutable
# =============================================================================

# Lock the audit configuration (prevents runtime changes)
# WARNING: After this, you must reboot to modify audit rules
# Uncomment only in production environments
# -e 2
