# docs/eoa/persistence_agent/scorecard.json
{
  "agent": "persistence",
  "version": "v2",
  "label": "PersistenceGuard V2",
  "eoa_date": "2026-02-16T17:04:09Z",
  "verdict": "PASS",
  "platform": "macOS (Apple Silicon)",
  "collector": "macOS LaunchAgents/Daemons, cron, shell profiles, SSH keys (live)",
  "duration_s": 120,
  "results": {
    "total_queue_rows": 12,
    "decode_ok": 12,
    "decode_fail": 0,
    "total_events": 16,
    "event_types": {
      "METRIC": 12,
      "SECURITY": 4
    },
    "mitre_techniques_observed": ["T1543", "T1546.004", "T1547"],
    "probes_declared": 8,
    "probes_fired": 2,
    "probes_fired_names": ["persistence_user_launch_agent", "persistence_shell_profile"],
    "empty_evidence_events": 0,
    "has_metric_data": 12,
    "has_security_event": 4,
    "has_attributes": 4
  },
  "resource_stability": {
    "rss_min_kb": 3936,
    "rss_max_kb": 52752,
    "rss_steady_state_kb": 39232,
    "cpu_max_pct": 0.0,
    "unexpected_tracebacks": 0
  },
  "trigger_run": {
    "actions": [
      "Created com.eoa.test.plist in ~/Library/LaunchAgents → detected CREATE",
      "Deleted com.eoa.test.plist → detected DELETE",
      "Modified .zshrc → detected shell profile change"
    ]
  },
  "fixes_applied": [
    "Implemented real macOS collector (LaunchAgents/Daemons, cron, shell profiles, SSH keys)",
    "Fixed all API mismatches with HardenedAgentBase",
    "setup(): cert check changed to warning-only",
    "collect_data(): SecurityEvent sub-messages + attributes map + heartbeat METRIC"
  ]
}
